# deploy/flowtragent-capture.service
[Unit]
Description=FlowTragent Live Capture Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flowtragent
Group=flowtragent
WorkingDirectory=/opt/FlowTragent
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/opt/FlowTragent/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=FLOWTRAGENT_LIVE_INCOMING_DIR=data/live/incoming
Environment=FLOWTRAGENT_SEGMENT_SECONDS=60
Environment=FLOWTRAGENT_PACKET_COUNT=0
Environment=FLOWTRAGENT_CAPTURE_PROFILE=balanced
Environment=FLOWTRAGENT_PYTHON=/opt/FlowTragent/.venv/bin/python
EnvironmentFile=-/etc/flowtragent/flowtragent.env
ExecStart=/usr/bin/env bash -c 'exec "$FLOWTRAGENT_PYTHON" /opt/FlowTragent/scripts/live_capture_worker.py --interface "$FLOWTRAGENT_INTERFACE" --output-dir "$FLOWTRAGENT_LIVE_INCOMING_DIR" --segment-seconds "$FLOWTRAGENT_SEGMENT_SECONDS" --packet-count "$FLOWTRAGENT_PACKET_COUNT" --profile "$FLOWTRAGENT_CAPTURE_PROFILE"'
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/FlowTragent/reports /opt/FlowTragent/data

[Install]
WantedBy=multi-user.target
